# async-server.service
[Unit]
Description=Async TCP/UDP Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/async_server 8080
Restart=always
RestartSec=5
User=nobody
Group=nogroup

[Install]
WantedBy=multi-user.target
